# antivirusInstall/daemon.service
[Unit]
Description=Najlepszy antyvirus 

[Service]
Type=simple
ExecStart=/usr/antivirus/daemon

[Install]
WantedBy=multi-user.target
